# deploy/supervisord.conf
; Runs the API and one or more workers in a single RunPod container.
; For multi-GPU pods, raise `numprocs` on the worker program (one worker per GPU)
; and set AIV_COMFY_ENDPOINTS accordingly.
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:api]
command=uvicorn app.main:app --host 0.0.0.0 --port 8000
directory=/workspace/api
autostart=true
autorestart=true
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:worker]
command=python -m app.workers.worker
directory=/workspace/api
autostart=true
autorestart=true
numprocs=1
process_name=%(program_name)s_%(process_num)02d
stopsignal=TERM
stopwaitsecs=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
